# etc/root/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/etc/init.d/rcK
::respawn:/usr/bin/getty -L console 0 vt100
